# templates/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
priority=10
directory=/opt/apps/solr-4.6.0/example
command=/usr/bin/java -jar start.jar
user=root
autostart=true
autorestart=true